# devel/etc/catkin/profile.d/cmake.lock
/home/govind/ariac_ws/devel/.private/nist_gear/etc/catkin/profile.d/cmake.lock